# supervisor.d/create_delivery_order_tmp_consumer.ini
[program:create_delivery_order_consumer_tmp]
command         = docker exec -i order-service ./ordersrv consumer create-delivery-order-tmp create-delivery-order-tmp create-do-group 0 0
autostart       = true
autorestart     = true
startsecs       = 0
user            = root
numprocs        = 1
stderr_logfile=/var/log/supervisor/create_delivery_order_consumer_tmp_err.log
stdout_logfile=/var/log/supervisor/create_delivery_order_consumer_tmp_out.log