# supervisor.d/delete_delivery_order_consumer.ini
[program:delete_delivery_order_consumer]
command         = docker exec -i order-service ./ordersrv consumer delete-delivery-order delete-delivery-order update-do-group 0 0
autostart       = true
autorestart     = true
startsecs       = 0
user            = root
numprocs        = 1
stderr_logfile=/var/log/supervisor/delete_delivery_order_consumer_err.log
stdout_logfile=/var/log/supervisor/delete_delivery_order_consumer_out.log