# supervisor.d/update_delivery_order_consumer.ini
[program:update_delivery_order_consumer]
command         = docker exec -i order-service ./ordersrv consumer update-delivery-order update-delivery-order update-do-group 0 0
autostart       = true
autorestart     = true
startsecs       = 0
user            = root
numprocs        = 1
stderr_logfile=/var/log/supervisor/update_delivery_order_consumer_err.log
stdout_logfile=/var/log/supervisor/update_delivery_order_consumer_out.log