# supervisor.d/update_sales_order_consumer.ini
[program:update_sales_order_consumer]
command         = docker exec -i order-service ./ordersrv consumer update-sales-order update-sales-order update-so-group 0 0
autostart       = true
autorestart     = true
startsecs       = 0
user            = root
numprocs        = 1
stderr_logfile=/var/log/supervisor/update_sales_order_consumer_err.log
stdout_logfile=/var/log/supervisor/update_sales_order_consumer_out.log