# supervisor.d/upload_so_file_consumer.ini
[program:upload_so_file_consumer]
command         = docker exec -i order-service ./ordersrv consumer upload-so-file upload-so-file upload-so-group 0 0
autostart       = true
autorestart     = true
startsecs       = 0
user            = root
numprocs        = 1
stderr_logfile=/var/log/supervisor/upload_so_file_consumer_err.log
stdout_logfile=/var/log/supervisor/upload_so_file_consumer_out.log