# supervisor.d/upload_do_file_consumer.ini
[program:upload_do_file_consumer]
command         = docker exec -i order-service ./ordersrv consumer upload-do-file upload-do-file upload-do-group 0 0
autostart       = true
autorestart     = true
startsecs       = 0
user            = root
numprocs        = 1
stderr_logfile=/var/log/supervisor/upload_do_file_consumer_err.log
stdout_logfile=/var/log/supervisor/upload_do_file_consumer_out.log